# dell_undock/etc_udev_rules.d_85-Latitude-E7450.rules
SUBSYSTEM=="usb", ENV{ID_MODEL_ID}=="2513", ENV{ID_VENDOR_ID}=="413c" RUN+="/etc/dell/dock.sh $env{ACTION} test"

