# dell_undock/etc_udev_rules.d_70-monitor.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/etc/dell/dock.sh"

